# ex/apds-gesture-ui.service
[Unit]
Description=APDS-9960 Gesture UI Service

[Service]
Type=simple
ExecStart=/usr/sbin/apds-gesture-ui
