# ex/apds-gesture-ui-sub.service
[Unit]
Description=APDS-9960 Gesture UI Service (MQTT Subscriber)

[Service]
Type=simple
ExecStart=/usr/sbin/apds-gesture-ui-sub
